# os/mkosi.extra/etc/systemd/system/getty@tty1.service.d/autologin.conf
# thriveOS is a single-purpose appliance: the local console belongs to whoever
# has physical access, so auto-login root there (no password) — you get a shell
# the moment it boots, which is the "headless but I can still see a terminal"
# behaviour we want. Remote access stays gated (sshd defaults). Drop this file to
# go back to a normal login: prompt.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
